# hostapd-init.service
[Unit]
Description=hostapd-init-service
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/scripts/
ExecStart=/bin/bash /root/scripts/hostapd-init.sh
SyslogIdentifier=hostpad-init

[Install]
WantedBy=multi-user.target


